# FreeBSD9-orig/etc/rc.d/newsyslog
#!/bin/sh
#
# $FreeBSD: release/9.0.0/etc/rc.d/newsyslog 197947 2009-10-10 22:17:03Z dougb $
#

# PROVIDE: newsyslog
# REQUIRE: cleanvar mountcritremote

. /etc/rc.subr

name="newsyslog"
rcvar=`set_rcvar`
required_files="/etc/newsyslog.conf"
command="/usr/sbin/${name}"
start_cmd="newsyslog_start"
stop_cmd=":"

newsyslog_start()
{
	check_startmsgs && echo -n 'Creating and/or trimming log files'
	${command} ${rc_flags}
	check_startmsgs && echo '.'
}

load_rc_config $name
run_rc_command "$1"
